# scripts/astronomical-events.service
# Systemd service for Astronomical Events Skill
# Install: sudo cp scripts/astronomical-events.service /etc/systemd/system/
#           sudo systemctl daemon-reload
#           sudo systemctl enable --now astronomical-events

[Unit]
Description=Astronomical Events Notification Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/urtzai/.openclaw/skills/astronomical-events/scripts/main.py schedule
WorkingDirectory=/home/urtzai/.openclaw/skills/astronomical-events
EnvironmentFile=/home/urtzai/.openclaw/skills/astronomical-events/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=astronomical-events

# Restart policy
Restart=on-failure
RestartSec=30
StartLimitBurst=5
StartLimitIntervalSec=300

# Security hardening
NoNewPrivileges=true
ProtectHome=read-only
ReadWritePaths=/home/urtzai/.openclaw/skills/astronomical-events/data

[Install]
WantedBy=multi-user.target
